# scheduler/cron.d/sample_tasks.cron
# BEGIN sample_tasks.cron
# TSWF sample cron entries
# Each line will be merged into user crontab by bin/install_cron.sh
# Format: <cron_time> cd <repo_path> && <command> # TSWF <task_name>

# Run Task A every 2 minutes
*/2 * * * * cd /path/to/tswf && ./cli/tswf.sh task run --name sampleA  # TSWF sampleA

# Run Task B every 3 minutes
*/3 * * * * cd /path/to/tswf && ./cli/tswf.sh task run --name sampleB  # TSWF sampleB
# END sample_tasks.cron
