# scheduler/cron.d/tswf.cron
*/2 * * * * cd /Users/user/Documents/GitHub/SOFE3200-TASK-SCHDULER && ./cli/tswf.sh task run --name sampleA # TSWF name=sampleA
*/3 * * * * cd /Users/user/Documents/GitHub/SOFE3200-TASK-SCHDULER && ./cli/tswf.sh task run --name sampleB # TSWF name=sampleB
